# supervisord.conf
[supervisord]
nodaemon=true

[program:opsmanager]
command=/startup.sh
